# deploy/production/systemd/chamados-ti-backup.timer
[Unit]
Description=Backup diário do Chamados TI

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=10m
Persistent=true
Unit=chamados-ti-backup.service

[Install]
WantedBy=timers.target
